# contrib/vigil.service
[Unit]
Description=Vigil — AI-powered Linux system guardian
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=/usr/local/bin/vigil
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vigil

[Install]
WantedBy=default.target
