# contrib/vigil.service
[Unit]
Description=Vigil — AI-powered Linux system guardian
After=network.target

[Service]
Type=simple
# User= is set by install.sh to the installing user
ExecStart=/usr/local/bin/vigil
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vigil

[Install]
WantedBy=default.target
